# webapp.path
[Unit]
Description=Watch for /opt/setDataBase.sh

[Path]
PathExists=/opt/setDataBase.sh

[Install]
WantedBy=multi-user.target